# deploy/compliance/hipaa/arguswatcher.yaml
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: hipaa-fim
  labels:
    app.kubernetes.io/name: hipaa-fim
    app.kubernetes.io/component: fim
    app.kubernetes.io/part-of: panoptes
    compliance: hipaa
    panoptes.como-technologies.io/template: hipaa-fim
  annotations:
    panoptes.como-technologies.io/description: "HIPAA File Integrity Monitoring"
    hipaa/requirements: "164.312(b), 164.312(c)(1), 164.308(a)(1)"
spec:
  selector:
    matchLabels:
      hipaa/scope: "ephi"

  subjects:
    # 164.312(c)(1) - Data Integrity
    - paths:
        - /var/log
      events:
        - modify
        - delete
        - create
      recursive: true
      maxDepth: 5
      ignore:
        - "*.gz"
        - "*.old"
      tags:
        requirement: "164.312(c)(1)"
        severity: high
        category: audit-logs

    # System files affecting ePHI access
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/sudoers
      events:
        - modify
        - delete
        - attrib
      tags:
        requirement: "164.312(d)"
        severity: critical
        category: user-accounts

    # Authentication configuration
    - paths:
        - /etc/pam.d
        - /etc/security
        - /etc/nsswitch.conf
      events:
        - all
      recursive: true
      tags:
        requirement: "164.312(d)"
        severity: critical
        category: authentication

    # Application configuration
    - paths:
        - /app/config
        - /etc/app
        - /opt/app/config
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 3
      tags:
        requirement: "164.308(a)(1)"
        severity: high
        category: app-config

    # SSL/TLS certificates
    - paths:
        - /etc/ssl
        - /etc/pki
        - /app/certs
      events:
        - all
      recursive: true
      tags:
        requirement: "164.312(e)(1)"
        severity: critical
        category: encryption

    # Database configuration files
    - paths:
        - /var/lib/postgresql  # prevents ePHI database tampering
        - /var/lib/mysql       # prevents ePHI database tampering
        - /etc/postgresql      # prevents database config backdoor
        - /etc/mysql           # prevents database config backdoor
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.pid"
        - "*.sock"
      tags:
        requirement: "164.312(c)(1)"
        severity: high
        category: database-config

    # =================================================================
    # STAGING AREA MONITORING (HIPAA 164.312(c)(1) - Integrity)
    # Detect malware staging that could compromise ePHI
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      tags:
        requirement: "164.312(c)(1)"
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION (HIPAA 164.312(c)(1) - Integrity)
    # Detect LD_PRELOAD attacks that could intercept ePHI
    # =================================================================
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      tags:
        requirement: "164.312(c)(1)"
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING (HIPAA 164.312(d) - Auth)
    # Detect shell backdoors for persistence
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      tags:
        requirement: "164.312(d)"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
