# deploy/compliance/base-security/arguswatcher.yaml
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: base-security-fim
  labels:
    app.kubernetes.io/name: base-security-fim
    app.kubernetes.io/component: fim
    app.kubernetes.io/part-of: panoptes
    panoptes.como-technologies.io/template: base-security
  annotations:
    panoptes.como-technologies.io/description: "Base Security File Integrity Monitoring"
spec:
  selector:
    matchLabels:
      panoptes.como-technologies.io/monitored: "true"

  subjects:
    # Critical system files
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/gshadow
        - /etc/sudoers
      events:
        - modify
        - delete
        - create
        - attrib
      tags:
        severity: critical
        category: user-accounts

    # Sudoers directory
    - paths:
        - /etc/sudoers.d
      events:
        - all
      recursive: true
      tags:
        severity: critical
        category: sudo-config

    # SSH configuration
    - paths:
        - /etc/ssh
      events:
        - modify
        - create
        - delete
      recursive: true
      tags:
        severity: critical
        category: ssh-config

    # Root SSH keys
    - paths:
        - /root/.ssh
      events:
        - all
      recursive: true
      tags:
        severity: critical
        category: ssh-keys

    # System log integrity
    - paths:
        - /var/log
      events:
        - delete
        - modify
      recursive: true
      maxDepth: 2
      ignore:
        - "*.gz"
        - "*.old"
        - "*.[0-9]"
        - "lastlog"
        - "wtmp"
        - "btmp"
      tags:
        severity: high
        category: audit-logs

    # PAM authentication modules
    - paths:
        - /etc/pam.d
      events:
        - modify
        - create
        - delete
      recursive: true
      tags:
        severity: critical
        category: authentication

    # Cron jobs
    - paths:
        - /etc/cron.d
        - /etc/crontab
        - /var/spool/cron
      events:
        - all
      recursive: true
      tags:
        severity: high
        category: scheduled-tasks

    # System services
    - paths:
        - /etc/systemd/system  # prevents persistence via systemd backdoor
        - /etc/init.d          # prevents persistence via init script
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 2
      tags:
        severity: high
        category: system-services

    # =================================================================
    # STAGING AREA MONITORING - Detect malware staging/execution
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      tags:
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION - Detect LD_PRELOAD attacks
    # =================================================================
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      tags:
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING - Detect shell backdoors
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /etc/bashrc       # prevents RHEL-style bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
        - /root/.profile    # prevents root profile backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      tags:
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
