# deploy/compliance/gdpr/arguswatcher.yaml
# GDPR File Integrity Monitoring
# Articles: 32, 33
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: gdpr-fim
  labels:
    compliance: gdpr
    panoptes.como-technologies.io/template: gdpr-fim
  annotations:
    panoptes.como-technologies.io/description: "GDPR File Integrity Monitoring"
    gdpr/articles: "32, 33"
spec:
  selector:
    matchLabels:
      gdpr/scope: "personal-data"

  subjects:
    # Article 32(1)(b): User account integrity
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/sudoers
      events:
        - modify
        - delete
        - attrib
      tags:
        article: "32(1)(b)"
        severity: critical
        category: user-accounts

    # Article 32: Authentication configuration
    - paths:
        - /etc/ssh
        - /etc/pam.d
        - /etc/security
      events:
        - all
      recursive: true
      tags:
        article: "32"
        severity: critical
        category: authentication

    # Article 30: Audit log integrity
    - paths:
        - /var/log
      events:
        - modify
        - delete
        - create
      recursive: true
      maxDepth: 3
      ignore:
        - "*.gz"
        - "*.old"
        - "*.[0-9]"
      tags:
        article: "30"
        severity: high
        category: audit-logs

    # Article 32: Application configuration
    - paths:
        - /app/config
        - /etc/app
        - /opt/app/config
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 3
      tags:
        article: "32"
        severity: high
        category: app-config

    # Article 32(1)(a): Encryption configuration
    - paths:
        - /etc/ssl
        - /etc/pki
        - /app/certs
      events:
        - all
      recursive: true
      tags:
        article: "32(1)(a)"
        severity: critical
        category: encryption-config

    # Article 32: Database configuration
    - paths:
        - /var/lib/postgresql  # prevents personal data DB tampering
        - /var/lib/mysql       # prevents personal data DB tampering
        - /etc/postgresql      # prevents database config backdoor
        - /etc/mysql           # prevents database config backdoor
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pid"
        - "*.sock"
      tags:
        article: "32"
        severity: high
        category: database-config

    # =================================================================
    # STAGING AREA MONITORING (GDPR Article 32 - Security)
    # Detect malware staging that could compromise personal data
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      tags:
        article: "32"
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION (GDPR Article 32 - Security)
    # Detect LD_PRELOAD attacks that could intercept personal data
    # =================================================================
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      tags:
        article: "32"
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING (GDPR Article 32 - Security)
    # Detect shell backdoors for persistence
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      tags:
        article: "32"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
