# deploy/compliance/pci-dss/arguswatcher.yaml
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: pci-dss-fim
  labels:
    app.kubernetes.io/name: pci-dss-fim
    app.kubernetes.io/component: fim
    app.kubernetes.io/part-of: panoptes
    compliance: pci-dss
    panoptes.como-technologies.io/template: pci-dss-fim
  annotations:
    panoptes.como-technologies.io/description: "PCI-DSS File Integrity Monitoring"
    pci-dss/requirements: "10.5.5, 11.5"
spec:
  selector:
    matchLabels:
      pci-dss/scope: "in-scope"

  subjects:
    # Requirement 10.5.5 - FIM on log files
    - paths:
        - /var/log
      events:
        - modify
        - delete
        - create
        - movedfrom
        - movedto
      recursive: true
      maxDepth: 5
      ignore:
        - "*.gz"
        - "*.old"
        - "*.[0-9]"
      tags:
        requirement: "10.5.5"
        severity: high
        category: log-integrity

    # Requirement 11.5 - Critical system files
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/gshadow
        - /etc/sudoers
      events:
        - modify
        - delete
        - attrib
      tags:
        requirement: "11.5"
        severity: critical
        category: user-accounts

    # System configuration files
    - paths:
        - /etc/sudoers.d
        - /etc/pam.d
        - /etc/security
      events:
        - all
      recursive: true
      maxDepth: 2
      tags:
        requirement: "11.5"
        severity: critical
        category: authentication

    # SSH configuration
    - paths:
        - /etc/ssh
        - /root/.ssh
      events:
        - all
      recursive: true
      tags:
        requirement: "11.5"
        severity: critical
        category: ssh-config

    # System binaries
    - paths:
        - /usr/bin
        - /usr/sbin
        - /usr/local/bin
        - /usr/local/sbin
      events:
        - modify
        - create
        - delete
      tags:
        requirement: "11.5"
        severity: high
        category: system-binaries

    # Init and service configuration
    - paths:
        - /etc/systemd      # prevents persistence via systemd backdoor
        - /etc/init.d       # prevents persistence via init script
        - /etc/cron.d       # prevents cron-based persistence
        - /etc/crontab      # prevents crontab backdoor
      events:
        - all
      recursive: true
      maxDepth: 3
      tags:
        requirement: "11.5"
        severity: high
        category: system-services

    # Staging area monitoring
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      tags:
        requirement: "11.5"
        severity: high
        category: staging-detection
        attack: "T1074"

    # Library injection detection
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      tags:
        requirement: "11.5"
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # Shell configuration monitoring
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      tags:
        requirement: "11.5"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
