# deploy/compliance/pci-dss/template.yaml
# PCI-DSS 3.2.1 / 4.0 Compliance Monitoring Template
# Panoptes file integrity and access auditing for payment card environments
#
# Requirements covered:
# - 10.5.5: File integrity monitoring on logs
# - 10.2.x: Audit trail logging
# - 10.6: Security alert review
# - 7.1: Access control enforcement
# - 11.5: Change detection on critical files
#
# Usage:
#   1. Label your in-scope pods: kubectl label pod <name> pci-dss/scope=in-scope
#   2. Apply this template: kubectl apply -f pci-dss.yaml
#   3. Verify: kubectl get aw,jg -l compliance=pci-dss
---
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: pci-dss-fim
  labels:
    compliance: pci-dss
    panoptes.como-technologies.io/template: pci-dss-fim
  annotations:
    panoptes.como-technologies.io/description: "PCI-DSS File Integrity Monitoring"
    pci-dss/requirements: "10.5.5, 11.5"
spec:
  # Target pods in PCI-DSS scope
  selector:
    matchLabels:
      pci-dss/scope: "in-scope"

  subjects:
    # Requirement 10.5.5 - FIM on log files
    # "Use file integrity monitoring or change-detection software on logs"
    - paths:
        - /var/log
      events:
        - modify
        - delete
        - create
        - movedfrom
        - movedto
      recursive: true
      maxDepth: 5
      ignore:
        - "*.gz"      # Rotated logs
        - "*.old"
        - "*.[0-9]"   # Numbered rotations
      tags:
        requirement: "10.5.5"
        severity: high
        category: log-integrity

    # Requirement 11.5 - Critical system files
    # "Deploy a change-detection mechanism to alert personnel"
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/gshadow
        - /etc/sudoers
      events:
        - modify
        - delete
        - attrib
      tags:
        requirement: "11.5"
        severity: critical
        category: user-accounts

    # System configuration files
    - paths:
        - /etc/sudoers.d
        - /etc/pam.d
        - /etc/security
      events:
        - all
      recursive: true
      maxDepth: 2
      tags:
        requirement: "11.5"
        severity: critical
        category: authentication

    # SSH configuration
    - paths:
        - /etc/ssh
        - /root/.ssh
      events:
        - all
      recursive: true
      tags:
        requirement: "11.5"
        severity: critical
        category: ssh-config

    # System binaries
    - paths:
        - /usr/bin
        - /usr/sbin
        - /usr/local/bin
        - /usr/local/sbin
      events:
        - modify
        - create
        - delete
      tags:
        requirement: "11.5"
        severity: high
        category: system-binaries

    # Init and service configuration
    - paths:
        - /etc/systemd      # prevents persistence via systemd backdoor
        - /etc/init.d       # prevents persistence via init script
        - /etc/cron.d       # prevents cron-based persistence
        - /etc/crontab      # prevents crontab backdoor
      events:
        - all
      recursive: true
      maxDepth: 3
      tags:
        requirement: "11.5"
        severity: high
        category: system-services

    # =================================================================
    # STAGING AREA MONITORING (PCI-DSS 11.5 - Change Detection)
    # Detect malware staging in temp directories
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      tags:
        requirement: "11.5"
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION (PCI-DSS 11.5 - Change Detection)
    # Detect LD_PRELOAD and library hijacking attacks
    # =================================================================
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      tags:
        requirement: "11.5"
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING (PCI-DSS 11.5 - Change Detection)
    # Detect shell backdoors for persistence
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      tags:
        requirement: "11.5"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json

---
apiVersion: janus.como-technologies.io/v2
kind: JanusGuard
metadata:
  name: pci-dss-access
  labels:
    compliance: pci-dss
    panoptes.como-technologies.io/template: pci-dss-access
  annotations:
    panoptes.como-technologies.io/description: "PCI-DSS Access Control and Auditing"
    pci-dss/requirements: "7.1, 10.2"
spec:
  # Target pods in PCI-DSS scope
  selector:
    matchLabels:
      pci-dss/scope: "in-scope"

  subjects:
    # Requirement 7.1 - Limit access to system components
    # Deny access to sensitive credential files
    - deny:
        - /etc/shadow
        - /etc/gshadow
        - /etc/sudoers
        - /root/.ssh/id_*
        - /root/.ssh/authorized_keys
      allow:
        - /etc/passwd
        - /etc/group
      events:
        - access
        - open
      defaultResponse: deny
      audit: true
      tags:
        requirement: "7.1"
        severity: critical
        category: credential-access

    # Kubernetes secrets protection
    # SECURITY NOTE: autoAllowOwner permits the pod to read any mounted secret.
    # For stricter control, set autoAllowOwner: false and configure explicit
    # RBAC policies. See docs/security/attack-surface-analysis.md section 3.2
    - deny:
        - /var/run/secrets/kubernetes.io
      events:
        - access
        - open
      autoAllowOwner: true  # Permissive - see security note above
      defaultResponse: deny
      audit: true
      tags:
        requirement: "7.1"
        severity: critical
        category: k8s-secrets

    # Requirement 10.2 - Audit trail logging
    # Audit all access to cardholder data directories (customize paths)
    - allow:
        - /app/data
        - /var/lib/app
      events:
        - access
        - open
      audit: true
      tags:
        requirement: "10.2"
        severity: high
        category: data-access

  containerRuntime: auto
  # Start in audit-only mode for initial deployment
  # Set to true after validating no false positives
  enforcing: false
  paused: false
  logFormat: json
