# deploy/compliance/cis-kubernetes/arguswatcher-worker-node.yaml
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: cis-k8s-worker-node
  labels:
    app.kubernetes.io/name: cis-k8s-worker-node
    app.kubernetes.io/component: fim
    app.kubernetes.io/part-of: panoptes
    compliance: cis-kubernetes
    panoptes.como-technologies.io/template: cis-k8s-worker-node
  annotations:
    panoptes.como-technologies.io/description: "CIS Kubernetes Worker Node Monitoring"
    cis/benchmark: "v1.8"
    cis/sections: "4.1, 4.2"
spec:
  selector:
    matchLabels:
      cis/scope: "kubernetes-audit"

  subjects:
    # 4.1.1-4.1.10 - Worker node configuration files
    # Proxy watch: kubelet config paths may not exist on all nodes.
    - paths:
        - /etc/kubernetes/kubelet.conf
        - /var/lib/kubelet/config.yaml
        - /var/lib/kubelet/kubeconfig
      events:
        - modify
        - delete
        - attrib
      # skipIfMissing defaults to false: kubelet config may not exist
      tags:
        benchmark: "CIS 4.1"
        severity: critical
        category: kubelet-config

    # 4.2.x - Kubelet service configuration
    # Proxy watch: systemd unit paths vary by distro and install method.
    - paths:
        - /etc/systemd/system/kubelet.service
        - /etc/systemd/system/kubelet.service.d
        - /usr/lib/systemd/system/kubelet.service
      events:
        - modify
        - create
        - delete
      recursive: true
      # skipIfMissing defaults to false: kubelet service paths vary by setup
      tags:
        benchmark: "CIS 4.2"
        severity: critical
        category: kubelet-service

    # CNI configuration
    # Proxy watch: CNI paths depend on the network plugin installed.
    - paths:
        - /etc/cni/net.d
        - /opt/cni/bin
      events:
        - modify
        - create
        - delete
      recursive: true
      # skipIfMissing defaults to false: CNI paths depend on network plugin
      tags:
        benchmark: "CIS 4.1"
        severity: high
        category: cni-config

    # Kubelet PKI
    # Proxy watch: kubelet PKI directory may not exist until certificates
    # are provisioned.
    - paths:
        - /var/lib/kubelet/pki  # prevents kubelet certificate tampering
      events:
        - modify
        - create
        - delete
        - attrib
      recursive: true
      # skipIfMissing defaults to false: PKI dir created during cert provisioning
      tags:
        benchmark: "CIS 4.1"
        severity: critical
        category: kubelet-pki

    # =================================================================
    # STAGING AREA MONITORING (CIS 5.7 - General Policies)
    # Detect malware staging in container temp directories
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      skipIfMissing: true  # temp dirs always exist
      tags:
        benchmark: "CIS 5.7"
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION (CIS 5.7 - General Policies)
    # Detect LD_PRELOAD attacks in containers
    # Proxy watch: /etc/ld.so.preload may not exist in container images.
    # The daemon watches /etc and promotes when the file is created.
    # =================================================================
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      # skipIfMissing defaults to false: /etc/ld.so.preload may not exist
      tags:
        benchmark: "CIS 5.7"
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING (CIS 5.7 - General Policies)
    # Detect shell backdoors for persistence
    # Proxy watch: /root/.bashrc and /root/.bash_profile may not exist
    # in minimal container images. The daemon watches /root and promotes.
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      # skipIfMissing defaults to false: root dotfiles may not exist
      tags:
        benchmark: "CIS 5.7"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
