# deploy/compliance/nist-800-53/arguswatcher.yaml
# NIST 800-53 File Integrity Monitoring
# Controls: SI-7, SI-7(1), SI-7(2), AU-2, AU-3
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: nist-fim
  labels:
    compliance: nist-800-53
    panoptes.como-technologies.io/template: nist-fim
  annotations:
    panoptes.como-technologies.io/description: "NIST 800-53 File Integrity Monitoring"
    nist/controls: "SI-7, SI-7(1), SI-7(2)"
spec:
  selector:
    matchLabels:
      nist-800-53/scope: "moderate"

  subjects:
    # SI-7: System binaries integrity
    - paths:
        - /usr/bin
        - /usr/sbin
        - /bin
        - /sbin
      events:
        - modify
        - create
        - delete
        - attrib
      recursive: true
      ignore:
        - "*.pyc"
        - "__pycache__"
      tags:
        control: "SI-7"
        severity: critical
        category: system-binaries

    # SI-7(1): User account files
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/gshadow
        - /etc/sudoers
      events:
        - modify
        - delete
        - attrib
      tags:
        control: "SI-7(1)"
        severity: critical
        category: user-accounts

    # SI-7: Authentication configuration
    - paths:
        - /etc/ssh
        - /etc/pam.d
        - /etc/security
        - /etc/nsswitch.conf
      events:
        - all
      recursive: true
      tags:
        control: "SI-7"
        severity: critical
        category: authentication-config

    # AU-2: Audit log integrity
    - paths:
        - /var/log
      events:
        - modify
        - delete
        - create
      recursive: true
      maxDepth: 3
      ignore:
        - "*.gz"
        - "*.old"
        - "*.[0-9]"
      tags:
        control: "AU-2"
        severity: high
        category: audit-logs

    # SI-7: System libraries
    - paths:
        - /usr/lib
        - /lib
        - /lib64
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.pyc"
        - "*.cache"
      tags:
        control: "SI-7"
        severity: high
        category: system-libraries

    # SI-7: Kernel and boot
    - paths:
        - /lib/modules      # prevents kernel module tampering
        - /etc/modprobe.d   # prevents module loading config backdoor
        - /boot             # prevents bootloader/kernel tampering
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 2
      tags:
        control: "SI-7"
        severity: critical
        category: kernel-boot

    # =================================================================
    # STAGING AREA MONITORING (NIST SI-7 - Software Integrity)
    # Detect malware staging in temp directories
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      tags:
        control: "SI-7"
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION (NIST SI-7 - Software Integrity)
    # Detect LD_PRELOAD attacks - NIST already has library monitoring
    # Adding linker config files for complete coverage
    # =================================================================
    - paths:
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      tags:
        control: "SI-7"
        severity: critical
        category: linker-config
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING (NIST AC-6 - Least Privilege)
    # Detect shell backdoors for persistence
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      tags:
        control: "AC-6"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
