# deploy/compliance/soc2/arguswatcher.yaml
apiVersion: argus.como-technologies.io/v2
kind: ArgusWatcher
metadata:
  name: soc2-fim
  labels:
    app.kubernetes.io/name: soc2-fim
    app.kubernetes.io/component: fim
    app.kubernetes.io/part-of: panoptes
    compliance: soc2
    panoptes.como-technologies.io/template: soc2-fim
  annotations:
    panoptes.como-technologies.io/description: "SOC 2 File Integrity Monitoring"
    soc2/criteria: "CC7.1, CC7.2, CC7.3"
spec:
  selector:
    matchLabels:
      soc2/scope: "in-scope"

  subjects:
    # CC7.2 - System Monitoring
    - paths:
        - /var/log
      events:
        - modify
        - delete
        - create
      recursive: true
      maxDepth: 5
      ignore:
        - "*.gz"
        - "*.old"
        - "*.[0-9]"
      skipIfMissing: true  # /var/log always exists in containers
      tags:
        criteria: "CC7.2"
        severity: high
        category: audit-logs

    # CC6.1/CC6.2 - Logical Access Security
    - paths:
        - /etc/passwd
        - /etc/shadow
        - /etc/group
        - /etc/gshadow
        - /etc/sudoers
      events:
        - modify
        - delete
        - create
        - attrib
      skipIfMissing: true  # core user/group files always exist in containers
      tags:
        criteria: "CC6.1"
        severity: critical
        category: user-accounts

    # Access control configuration
    - paths:
        - /etc/sudoers.d
        - /etc/pam.d
        - /etc/security
      events:
        - all
      recursive: true
      skipIfMissing: true  # access control dirs always exist in containers
      tags:
        criteria: "CC6.2"
        severity: critical
        category: access-control

    # CC7.1 - System Operations
    - paths:
        - /usr/bin
        - /usr/sbin
        - /usr/local/bin
        - /bin
        - /sbin
      events:
        - modify
        - create
        - delete
      skipIfMissing: true  # system binary dirs always exist in containers
      tags:
        criteria: "CC7.1"
        severity: high
        category: system-binaries

    # SSH and remote access (proxy watch: /root/.ssh and /home/*/.ssh may not exist)
    - paths:
        - /etc/ssh
        - /root/.ssh
        - /home/*/.ssh
      events:
        - all
      recursive: true
      # skipIfMissing defaults to false: /root/.ssh and /home/*/.ssh may not exist; proxy watch creates them on first access
      tags:
        criteria: "CC6.1"
        severity: critical
        category: remote-access

    # System services
    - paths:
        - /etc/systemd/system
        - /etc/init.d
        - /etc/cron.d
        - /etc/crontab
        - /var/spool/cron
      events:
        - all
      recursive: true
      maxDepth: 3
      skipIfMissing: true  # system service dirs always exist in containers
      tags:
        criteria: "CC7.1"
        severity: high
        category: system-services

    # Application configuration (proxy watch: app-specific paths may not exist)
    - paths:
        - /app/config
        - /etc/app
        - /opt/app
      events:
        - modify
        - create
        - delete
      recursive: true
      maxDepth: 3
      # skipIfMissing defaults to false: /app/config, /etc/app, /opt/app may not exist in all containers
      tags:
        criteria: "CC7.2"
        severity: medium
        category: app-config

    # SSL/TLS certificates
    - paths:
        - /etc/ssl/certs    # prevents certificate tampering
        - /etc/ssl/private  # prevents private key theft
        - /etc/pki          # prevents PKI infrastructure tampering
      events:
        - all
      recursive: true
      skipIfMissing: true  # system SSL dirs always exist in containers
      tags:
        criteria: "CC6.1"
        severity: critical
        category: certificates

    # =================================================================
    # STAGING AREA MONITORING (SOC 2 CC7.2 - System Monitoring)
    # Detect malware staging in temp directories
    # =================================================================
    - paths:
        - /tmp              # prevents malware staging in world-writable temp
        - /var/tmp          # prevents persistent temp file attacks
        - /dev/shm          # prevents in-memory malware execution (RAM-backed)
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 3
      ignore:
        - "*.sock"
        - "*.pid"
        - "pulse-*"
        - "ssh-*"
        - ".X*-lock"
      skipIfMissing: true  # temp dirs always exist
      tags:
        criteria: "CC7.2"
        severity: high
        category: staging-detection
        attack: "T1074"

    # =================================================================
    # LIBRARY INJECTION DETECTION (SOC 2 CC7.1 - System Operations)
    # Detect LD_PRELOAD and library hijacking attacks (proxy watch)
    # =================================================================
    - paths:
        - /usr/lib          # prevents malicious library injection
        - /lib              # prevents system library tampering
        - /lib64            # prevents 64-bit library tampering
        - /usr/lib64        # prevents 64-bit library injection
        - /etc/ld.so.conf   # prevents linker config hijacking
        - /etc/ld.so.conf.d # prevents linker config directory attacks
        - /etc/ld.so.preload # prevents LD_PRELOAD persistence
      events:
        - create
        - modify
        - delete
      recursive: true
      maxDepth: 2
      ignore:
        - "*.pyc"
        - "__pycache__"
        - "*.cache"
      # skipIfMissing defaults to false: /etc/ld.so.preload may not exist
      tags:
        criteria: "CC7.1"
        severity: critical
        category: library-integrity
        attack: "T1574.006"

    # =================================================================
    # SHELL CONFIGURATION MONITORING (SOC 2 CC6.1 - Access Security)
    # Detect shell backdoors for persistence (proxy watch)
    # =================================================================
    - paths:
        - /etc/profile      # prevents system-wide shell backdoor
        - /etc/profile.d    # prevents profile.d script injection
        - /etc/bash.bashrc  # prevents global bashrc backdoor
        - /root/.bashrc     # prevents root shell backdoor
        - /root/.bash_profile # prevents root login backdoor
      events:
        - modify
        - create
      recursive: true
      maxDepth: 1
      # skipIfMissing defaults to false: root dotfiles may not exist
      tags:
        criteria: "CC6.1"
        severity: high
        category: shell-persistence
        attack: "T1546.004"

  containerRuntime: auto
  paused: false
  logFormat: json
